# net-vpn/windscribe-bin/files/windscribe-helper.initd
#!/sbin/openrc-run

name="Windscribe Helper"
description="Windscribe VPN Background Service"
command="/opt/windscribe/helper"
command_background="true"
pidfile="/run/windscribe-helper.pid"

depend() {
	use net
	after network
	after firewall
	use dns
}
